# com.oraclex.brain-mirror.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd: auto-mirror the mother brain every 6h. Replace HOME paths.
     Install: sed "s#HOME#$HOME#g" → ~/Library/LaunchAgents/ → launchctl load -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.oraclex.brain-mirror</string>
  <key>ProgramArguments</key><array>
    <string>/bin/bash</string><string>HOME/.oraclex/brain-mirror.sh</string>
  </array>
  <key>StartInterval</key><integer>21600</integer>
  <key>RunAtLoad</key><true/>
  <key>StandardOutPath</key><string>HOME/.oraclex/mirror.log</string>
  <key>StandardErrorPath</key><string>HOME/.oraclex/mirror.log</string>
</dict></plist>
